# newContainer/newContainer/var/lib/dpkg/info/mongodb-10gen.list
/.
/usr
/usr/bin
/usr/bin/mongoimport
/usr/bin/mongodump
/usr/bin/mongoexport
/usr/bin/mongos
/usr/bin/mongorestore
/usr/bin/mongo
/usr/bin/mongostat
/usr/bin/mongotop
/usr/bin/mongofiles
/usr/bin/mongod
/usr/bin/bsondump
/usr/sbin
/usr/share
/usr/share/doc
/usr/share/doc/mongodb-10gen
/usr/share/doc/mongodb-10gen/changelog.gz
/usr/share/doc/mongodb-10gen/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mongofiles.1.gz
/usr/share/man/man1/mongoimport.1.gz
/usr/share/man/man1/mongostat.1.gz
/usr/share/man/man1/mongodump.1.gz
/usr/share/man/man1/mongo.1.gz
/usr/share/man/man1/mongos.1.gz
/usr/share/man/man1/mongorestore.1.gz
/usr/share/man/man1/mongoexport.1.gz
/usr/share/man/man1/bsondump.1.gz
/usr/share/man/man1/mongosniff.1.gz
/usr/share/man/man1/mongod.1.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/mongodb-10gen
/var
/var/lib
/var/lib/mongodb
/etc
/etc/init.d
/etc/init.d/mongodb
/etc/mongodb.conf
